# docker_stuff/ssh/startup_machine.service
[Unit]
Description=Startup the system for CTF

[Service]
Type=simple
ExecStart=/root/startup.sh
